# imagor_docker/services/imagor/run
#!/command/with-contenv bash
exec imagor