# imagor_docker/services/cron/run
#!/command/with-contenv bash
exec supercronic /etc/crontab
